# systemd/pe-format2.service.in
[Unit]
Description=Re-compile pe-format2 config

[Service]
Type=oneshot
ExecStart=@sbindir@/pe-format2-setup
